# custom/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStartPre=/usr/local/bin/update-issue
